# config/includes.chroot/etc/xdg/autostart/ai-assistant.desktop
[Desktop Entry]
Name=AI Assistant
Exec=/usr/local/bin/ai_assistant.py
Terminal=false
Type=Application
X-GNOME-Autostart-enabled=true